# deploy/systemd/agent-reactor-server.service
# agent-reactor server — merged pty session daemon + HTTP/WS gateway in one
# process. Owns sessions, sockets, and on-disk state under -data-dir.
# Restarting this unit drops every attached browser tab; sessions are
# recreated from disk on the next boot.
#
# Default bind is 127.0.0.1:8443 + -insecure (loopback only). Operators
# exposing to LAN/internet should set -tls-cert / -tls-key and re-bind via
# server.env. Bearer-token auth is always on (see -token-file below) unless
# -no-auth is explicitly added to ExecStart (loopback-only guard still applies).

[Unit]
Description=agent-reactor server (session daemon + HTTP/WS gateway)
# Operator guide: docs/user/systemd.md in the agent-reactor source repo.

[Service]
Type=notify
NotifyAccess=main
StateDirectory=agent-reactor
# -data-dir tells the server where to write its socket, pid file, sessions
# snapshot and logs. -token-file persists a bearer token across restarts so
# bookmarked browser URLs survive a unit reload.
ExecStart=%h/.local/bin/agent-reactor-server \
  -addr 127.0.0.1:8443 -insecure \
  -data-dir %S/agent-reactor \
  -token-file %S/agent-reactor/server.token
# Operator override file (TLS, alternate addr, etc.). Leading `-` tolerates absence.
EnvironmentFile=-%S/agent-reactor/server.env

Restart=on-failure
RestartSec=2s
# srv.Shutdown(5s) in cmd/server's signal handler; 10s leaves headroom.
TimeoutStopSec=10s

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
